# mindx_observability/podman_quadlets/alertmanager.container
[Unit]
Description=Alertmanager 0.32.0 — mindX paging
After=network-online.target obs-network.service
Wants=network-online.target obs-network.service

[Container]
Image=quay.io/prometheus/alertmanager:v0.32.0
ContainerName=alertmanager
AutoUpdate=registry
Network=obs.network
Volume=%h/obs/alertmanager/alertmanager.yml:/etc/alertmanager/alertmanager.yml:Z,ro
Volume=%h/obs/alertmanager/smtp_app_password:/etc/alertmanager/smtp_app_password:Z,ro
Volume=%h/obs/alertmanager_data:/alertmanager:Z
PublishPort=127.0.0.1:9093:9093
Exec=--config.file=/etc/alertmanager/alertmanager.yml \
     --storage.path=/alertmanager \
     --web.external-url=https://alerts.pythai.net/
LogDriver=journald
NoNewPrivileges=true

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=300
OOMScoreAdjust=300

[Install]
WantedBy=default.target
