# mindx_observability/podman_quadlets/prometheus.container
[Unit]
Description=Prometheus 3.11.1 — mindX TSDB
After=network-online.target obs-network.service
Wants=network-online.target obs-network.service

[Container]
Image=quay.io/prometheus/prometheus:v3.11.1
ContainerName=prometheus
AutoUpdate=registry
Network=obs.network
Volume=%h/obs/prometheus/prometheus.yml:/etc/prometheus/prometheus.yml:Z,ro
Volume=%h/obs/prometheus/rules:/etc/prometheus/rules:Z,ro
Volume=%h/obs/prometheus_data:/prometheus:Z
PublishPort=127.0.0.1:9090:9090
Environment=GOGC=40
Exec=--config.file=/etc/prometheus/prometheus.yml \
     --storage.tsdb.path=/prometheus \
     --storage.tsdb.retention.time=3d \
     --storage.tsdb.retention.size=4GB \
     --web.enable-lifecycle \
     --web.enable-admin-api \
     --web.external-url=https://prom.pythai.net/ \
     --web.route-prefix=/
LogDriver=journald
NoNewPrivileges=true

[Service]
Restart=no
RestartSec=5
TimeoutStartSec=900
MemoryMax=256M
CPUQuota=20%
OOMScoreAdjust=500

# Phase 1.2: Prometheus is OFF BY DEFAULT.
# No [Install] section = no auto-start at boot.
# netdata is the always-on monitoring UI (https://netdata.pythai.net).
# Bring the full Prom stack up on demand via `scripts/prom_on.sh` for PromQL
# sessions or to arm email alerts. `scripts/prom_off.sh` shuts it back down,
# freeing ~416 MB resident.
